# cfg/system/profile.d/user.sh
# /etc/profile.d/user.sh
#
# This file contains global default env vars
#
listop prepend PROFILES /etc/profile.d/user.sh
